# recipes-bsp/acpi-tables/files/acpi-tables-load.service
[Unit]
Description=ACPI tables load service

[Service]
ExecStart=/usr/bin/acpi-tables-load
Restart=on-failure

[Install]
WantedBy=basic.target

